# mousemux.service
[Unit]
Description=Mouse-to-keyboard remapper daemon
After=systemd-udevd.service

[Service]
Type=simple
User=root
ExecStart=/usr/local/bin/mousemux --config /etc/mousemux/config.yaml
Restart=on-failure
RestartSec=2
NoNewPrivileges=yes
PrivateTmp=yes
ProtectHome=yes
ProtectSystem=strict
ReadOnlyPaths=/etc/mousemux
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
