# mousefold.service
[Unit]
Description=Mouse-to-keyboard remapper daemon
After=systemd-udevd.service

[Service]
Type=simple
User=root
ExecStart=/usr/local/bin/mousefold --config /etc/mousefold/config.yaml
ExecReload=/usr/local/bin/mousefold reload --config /etc/mousefold/config.yaml
Restart=on-failure
RestartSec=2
RuntimeDirectory=mousefold
NoNewPrivileges=yes
PrivateTmp=yes
ProtectHome=yes
ProtectSystem=strict
ReadOnlyPaths=/etc/mousefold
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
